# swbka/services/swbka.timer
[Unit]
Description=swbka timer

[Timer]
# 每天21点备份, 如有多个时间点, 复制该行
OnCalendar= *-*-* 21:00:00
Unit=swbka.service